# systemd/pulseaudio.service
[Unit]
Description=Pulseaudio

[Service]
Type=forking
ExecStart=/usr/bin/pulseaudio --start

[Install]
WantedBy=multi-user.target
